# task-sync.service
[Unit]
Description=auto sync taskwarrior

[Service]
ExecStart=/usr/bin/task sync    

[Install]
WantedBy=default.target
